# pd-mapper.service.in
[Unit]
Description=Qualcomm PD mapper service
ConditionDirectoryNotEmpty=/sys/class/remoteproc

[Service]
ExecStart=PD_MAPPER_PATH/pd-mapper
Restart=always

[Install]
WantedBy=multi-user.target
